# bioslow.service.j2
[Unit]
Description=Logging slow I/O ops

[Service]
ExecStart=/usr/local/bin/bioslow {{ bioslow_args }}
